# build_files/config/systemd/sddm.service.d/gaming-mode.conf
[Unit]
Conflicts=kodi-gbm.service
Wants=lactd.service
